# private_dot_config/systemd/user/delta-theme-startup.service
[Unit]
Description=Delta Theme Auto-setter on Startup
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/home/shinjungtae/.local/bin/delta-auto-theme.sh
RemainAfterExit=no

[Install]
WantedBy=default.target