# private_dot_config/systemd/user/hyprland-theme.timer
[Unit]
Description=Switch Hyprland Theme at scheduled times

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=false

[Install]
WantedBy=timers.target
