# private_dot_config/systemd/user/hyprdynamicmonitors.service
[Unit]
Description=HyprDynamicMonitors - Dynamic monitor configuration for Hyprland
Documentation=https://fiffeek.github.io/hyprdynamicmonitors/
After=hyprdynamicmonitors-prepare.service

[Service]
Type=exec
ExecStart=%h/.local/bin/hyprdynamicmonitors run --enable-lid-events
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
