# private_dot_config/systemd/user/hyprland-theme.service
[Unit]
Description=Apply Theme based on time (triggered by timer)
After=graphical-session.target

[Service]
Type=oneshot
KillMode=process
# Determine theme based on current hour
ExecStart=/bin/bash -c 'HOUR=$(date +%%H); if [ "$HOUR" -eq 0 ]; then %h/.local/bin/toggle-theme midnight system; fi'
