# private_dot_config/systemd/user/hyprland-theme-init.service
[Unit]
Description=Initialize Hyprland Theme Based on Time
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/hyprland-theme-init.sh
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
RemainAfterExit=no

[Install]
WantedBy=default.target
