# private_dot_config/systemd/user/opencode-server.service
[Unit]
Description=OpenCode Headless Server
Documentation=https://opencode.ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.opencode/bin/opencode serve --hostname 127.0.0.1 --port 42800
Restart=on-failure
RestartSec=10

# Environment
Environment=PATH=%h/.nvm/versions/node/v24.12.0/bin:%h/.opencode/bin:%h/.cargo/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=%h
Environment=OPENCODE_DISABLE_CLAUDE_CODE=1

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
